# contrib/openbsd/mitra
#!/bin/ksh

daemon="/usr/local/bin/mitra"
daemon_flags="server"
daemon_user="mitra"

. /etc/rc.d/rc.subr

rc_bg=YES

rc_cmd $1
